# etc/cron.d/ccnr3-update-toc
*/15 * * * * www-data /var/www/szen.in/www/n/app/share/cron/update-toc 60
